# deploy/vlab-frontend.service
[Unit]
Description=VLab Frontend (Next.js)
After=network.target vlab-backend.service

[Service]
Type=simple
User=vlab
Group=vlab

WorkingDirectory=/opt/vlab/frontend
ExecStart=/usr/bin/npm start
Restart=on-failure
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=vlab-frontend

[Install]
WantedBy=multi-user.target
